# google-startup-scripts/etc/init.d/google
#!/sbin/runscript
# Copyright 2013 Google Inc. All Rights Reserved.
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#     http://www.apache.org/licenses/LICENSE-2.0
#
#     Unless required by applicable law or agreed to in writing, software
#     distributed under the License is distributed on an "AS IS" BASIS,
#     WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
#     See the License for the specific language governing permissions and
#     limitations under the License.
#
### BEGIN INIT INFO
# Provides:    gce_onboot
# X-Start-Before:    ssh
# Required-Start:    $local_fs $network $syslog
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:
# Short-Description: Google Compute Engine on-boot services
# Description:       This launches the Google Compute Engine
#        VM initialization scripts.
### END INIT INFO

description="Google Compute Engine on-boot services"
NAME="onboot"

depend ()
{
  before ssh
  need localmount net
  use logger
  provide gce_onboot
}

#
# Function that starts the daemon/service
#
do_start()
{
  /usr/share/google/onboot
}

start () {
    ebegin "Starting $DESC $NAME"
    do_start
    case "$?" in
      0) eend 0 ;;
      *) eend 1 ;;
    esac
}

